# deploy/backup/cinta-dhuafa-backup.cron
30 2 * * * root /usr/local/bin/cinta-dhuafa-backup >> /var/log/cinta-dhuafa-backup.log 2>&1
